# supervisor/repo.conf
[program:finapp]
command=/personalwebsite gunicorn app:app -b localhost:8000
directory=/home/ubuntu/repo
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
